# server/supervisor/production.ini
[program:dataexpress-site-production]
directory = /home/devuser/webapps/dataexpress-site-env/dataexpress-site
command = /home/devuser/webapps/dataexpress-site-env/bin/uwsgi --ini /home/devuser/webapps/dataexpress-site-env/dataexpress-site/server/uwsgi/production.ini
autostart = true
autorestart = true
